# nextcloud/nextcloud-cron/crontabs/www-data
PHP_MEMORY_LIMIT=512M
*/5 * * * * php -f /var/www/html/cron.php
*/5 * * * * php /var/www/html/occ preview:pre-generate
